# platform/linux/systemd/qotd-fetch.timer
[Unit]
Description=Fetch new quote daily and on boot

[Timer]
OnBootSec=5min
OnCalendar=daily

[Install]
WantedBy=timers.target
