# platform/linux/systemd/qotdsh-fetch.service
[Unit]
Description=Fetch new quote

[Service]
Type=oneshot
ExecStart=/usr/local/bin/qotdsh fetch
